# dev/platformd/platformd.service
[Unit]
Requires=crio.service

[Service]
ExecStart=/usr/bin/platformd